# config/systemd/raffle-kiosk.service
[Unit]
Description=Raffle Entry Collector (Rails/Puma)
After=local-fs.target

[Service]
Type=simple
User=raffle
Group=raffle
WorkingDirectory=/opt/raffle
EnvironmentFile=/opt/raffle/.env
Environment=RAILS_ENV=production

ExecStart=/opt/raffle/.rbenv/shims/bundle exec puma -C config/puma.rb

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

# Memory safety net — restart if leak develops (Pi 4 has 4GB)
MemoryMax=512M

# Security: runs as unprivileged 'raffle' user with nologin shell.
# The Pi is a single-purpose kiosk — filesystem sandboxing (ProtectSystem,
# ProtectHome) caused SQLite open failures and isn't worth the complexity.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
