# config/systemd/raffle-backup.timer
[Unit]
Description=Raffle Entry Collector USB Backup Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
